# deploy/bur-web.service
# /etc/systemd/system/bur-web.service on 192.168.0.79
# Runs the BUR web app on 127.0.0.1:3000; nginx proxies :82 -> :3000.
[Unit]
Description=BUR — Blueberry User Repository (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bur
WorkingDirectory=/opt/bur/web
EnvironmentFile=/opt/bur/web/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/bur/web/.next
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
